# cwm.conf
[program:cwm]
command=/usr/local/bin/gunicorn wsgi --bind localhost:4000 -w 3  --worker-class gevent --log-level debug
directory=/var/www/cwm
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/cwm_stdout.log
stderr_logfile=/var/log/supervisor/cwm_stderr.log

# supervisor 中文 bug
environment=LANG="en_US.utf8", LC_ALL="en_US.UTF-8", LC_LANG="en_US.UTF-8"